# betclient/files/etc/init.d/betclient
#!/bin/sh /etc/rc.common
# iup at boot
# Copyright (C) 2007 OpenWrt.org

START=98
STOP=15
USE_PROCD=1

start_service() {
        betclient &
}
